# docker/supervisord.conf
[supervisord]
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/app/logs
nodaemon=true

[program:uvicorn]
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/app/logs/uvicorn.log
stderr_logfile=/app/logs/uvicorn.err.log

[program:nextjs]
command=/usr/bin/node server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/app/logs/next.log
stderr_logfile=/app/logs/next.err.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.err.log
